# hybrid-config/hybrid-send-logs.timer
[Unit]
Description=Send Hybrid's logs

[Timer]
OnCalendar=*-*-* 01:12:00
Persistent=true

[Install]
WantedBy=timers.target
